# node-package/launchd/ai.agora.autoupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.agora.autoupdate</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__AGORA_DIR__/scripts/auto_update.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>AGORA_DIR</key>
        <string>__AGORA_DIR__</string>
    </dict>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>__AGORA_DIR__/update.log</string>
    <key>StandardErrorPath</key>
    <string>__AGORA_DIR__/update.log</string>
</dict>
</plist>
